# scripts/com.polarisor.evolution-loop-cron.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polarisor.evolution-loop-cron</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>scripts/evolution-loop-cron-tick.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_POLARUI_REPO_ROOT</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>POLAR_EVOLUTION_LIVE_LLM</key>
        <string>0</string>
    </dict>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>REPLACE_WITH_POLARUI_REPO_ROOT/logs/evolution-loop-cron.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_POLARUI_REPO_ROOT/logs/evolution-loop-cron.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
